# templates/teleport.service.j2
[Unit]
Description=Teleport SSH Service

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/teleport start --config=/etc/teleport.yaml

[Install]
WantedBy=multi-user.target